# scripts/services/com.aws.profile-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aws.profile-bridge</string>
    
    <key>ProgramArguments</key>
    <array>
        <!-- Will be replaced with Python 3.12+ path -->
        <string>/usr/local/bin/python3.12</string>
        <string>-m</string>
        <string>aws_profile_bridge</string>
        <string>api</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>ENV</key>
        <string>production</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/.aws/logs/aws_profile_bridge_api.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/.aws/logs/aws_profile_bridge_api.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/USERNAME</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
